# launchd/com.fchangjun.gold-price-alert.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fchangjun.gold-price-alert</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/miniconda3/bin/python3</string>
        <string>/Users/baba/Desktop/黄金/gold_alert.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/baba/Desktop/黄金</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/baba/Desktop/黄金/logs/gold-price-alert.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/baba/Desktop/黄金/logs/gold-price-alert.error.log</string>
</dict>
</plist>
